# src/setup/launchd/com.clomia.automata.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clomia.automata.supervisor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/.local/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/USERNAME/dev/claude-automata</string>
        <string>python</string>
        <string>-m</string>
        <string>system.supervisor</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/dev/claude-automata</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ExitTimeOut</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/USERNAME/.local/bin:/Users/USERNAME/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>HOME</key>
        <string>/Users/USERNAME</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/USERNAME/dev/claude-automata/logs/supervisor-launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/dev/claude-automata/logs/supervisor-launchd-stderr.log</string>

    <key>Nice</key>
    <integer>5</integer>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>
